# systemd/px-alive.service
[Unit]
Description=PiCar-X Idle Alive Daemon (organic servo motion)
After=network.target px-mind.service px-api-server.service
StartLimitIntervalSec=0

[Service]
Type=simple
User=root
WorkingDirectory=/home/pi/picar-x-hacking
Environment=HOME=/home/pi
Environment=PX_ENABLE_ADULT_PERSONAS=0
ExecStart=/home/pi/picar-x-hacking/bin/px-alive --no-face
Restart=always
RestartSec=15
KillSignal=SIGTERM
TimeoutStopSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
