# systemd/px-frigate-stream.service
[Unit]
Description=PiCar-X Frigate camera stream (local go2rtc RTSP server, pull model)
After=network.target
# Architecture: go2rtc 1.9.7 runs locally with rpicam-vid as exec: source.
# Frigate on pi5-hailo pulls from rtsp://picar.local:8554/picar-x (pull model).
# This avoids Frigate-embedded go2rtc's RTSP ANNOUNCE bug (drops after 1 frame).
#
# px-alive runs with --no-face, so it does not open the libcamera pipeline and
# does not contend with rpicam-vid. No camera-yield handshake is needed.

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/picar-x-hacking
Environment=HOME=/home/pi
Environment=PATH=/home/pi/picar-x-hacking/.venv/bin:/usr/local/bin:/usr/bin:/bin
ExecStart=/home/pi/picar-x-hacking/bin/px-frigate-stream
Restart=on-failure
RestartSec=15
# Give rpicam-vid time to release the camera on stop
TimeoutStopSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
