# systemd/spark-pip-cleanup.service
[Unit]
Description=Sweep leftover pip scratch dirs in /tmp (failed install cleanup)

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/spark-pip-cleanup.sh
